# arm-docker/systemctl/rsr-serial-bridge.service
[Unit]
Description=RSR Serial Bridge Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/root/rsr
EnvironmentFile=-/root/rsr/.env
ExecStart=/usr/bin/python3 /root/rsr/main.py --config /root/rsr/settings.json
Restart=always
RestartSec=3
User=root

[Install]
WantedBy=multi-user.target
